# .apps.conf
[program:nginx]
user = root
command = /usr/sbin/nginx
